# ccc.sitebuilder.cron
* * * * * root /home/okoeroo/dvl/clib/streamcap/sitecreator.sh
